# crontab
*/1 * * * * root echo "Hello world" >> /var/log/cron.log
*/1 * * * * root echo "hey people" >> /var/log/words.log
